# ansible/files/lib/systemd/system/slide.rabbit-shocker.org.updater.timer
[Unit]
Description=Update slide.rabbit-shocker.org daily

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
